# apps/runner/packaging/deb/DEBIAN/prerm
#!/bin/bash
set -e

# Stop the service if it's running
if systemctl is-active --quiet boxlite-runner.service; then
    systemctl stop boxlite-runner.service
fi

# Disable the service
if systemctl is-enabled --quiet boxlite-runner.service; then
    systemctl disable boxlite-runner.service
fi

exit 0
